# root/etc/e-smith/templates/etc/crontab/durep

# Generate DUREP reports

01 0 * * * root /usr/sbin/durep.daily
